# app/code/App/RabbitMq/fixtures/supervisor.conf
[program:product_updater]
user=magento2
command= php /home/www/magento2/www/magento2.dev/bin/magento rabbitmq:worker:updater
process_name= product_updater
numprocs=1
autostart=true
autorestart=true
startretries=5
